# scripts/systemd/glitch-backend.service
[Unit]
Description=Glitch Lamp backend (uvicorn)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/yomope/glitch-lamp
Environment=HOME=/home/yomope UVICORN_HOST=0.0.0.0 UVICORN_PORT=18000 PYTHONUNBUFFERED=1
ExecStart=/home/yomope/glitch-lamp/venv/bin/uvicorn backend.main:app --host 0.0.0.0 --port 18000
User=yomope
Group=yomope
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
